# deploy/systemd/autocore-gateway.service
[Unit]
Description=AutoCore Gateway - MQTT Bridge
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=autocore
Group=autocore
WorkingDirectory=/opt/autocore/gateway

# Ambiente virtual isolado
Environment="PATH=/opt/autocore/gateway/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/autocore/database"

# Carregar variáveis de ambiente
EnvironmentFile=/opt/autocore/.env

# Comando com venv específico
ExecStart=/opt/autocore/gateway/.venv/bin/python src/main.py

# Restart policy
Restart=always
RestartSec=10

# Logs
StandardOutput=journal
StandardError=journal

# Limites de recursos (ajustar conforme necessário)
# MemoryLimit=256M
# CPUQuota=75%

[Install]
WantedBy=multi-user.target