# deploy/systemd/autocore-config-app.service
[Unit]
Description=AutoCore Config App - Web Interface
After=network.target autocore-gateway.service

[Service]
Type=simple
User=autocore
Group=autocore
WorkingDirectory=/opt/autocore/config-app/backend

# Ambiente virtual isolado
Environment="PATH=/opt/autocore/config-app/backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/autocore/database"

# Carregar variáveis de ambiente
EnvironmentFile=/opt/autocore/.env

# Comando com venv específico
ExecStart=/opt/autocore/config-app/backend/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8081 --workers 1

# Restart policy
Restart=always
RestartSec=10

# Logs
StandardOutput=journal
StandardError=journal

# Limites de recursos (ajustar conforme necessário)
# MemoryLimit=256M
# CPUQuota=75%

[Install]
WantedBy=multi-user.target